# rsync_server/rsync-exclude.list
/usr/local/nvidia
/usr/local/nvidia/*
/proc/*
/sys/*
/.autofsck
/boot/grub*/grub.cfg
/cm/local/apps/openldap/etc/certs/ldap.key
/cm/local/apps/openldap/etc/certs/ldap.pem
/data/*
/home/*
/local/*
/scratch/*
/tmp/*
/var/log/*
/var/tmp/*
/var/spool/mail/*
/var/spool/postfix/*
/var/lib/etcd/*
/var/lib/conatinerd/*
/var/lib/calico/*
/var/crash/*
/lib/modules/*/kernel/drivers/video/nvidia*.ko
/usr/lib/modules/*/kernel/drivers/video/nvidia*.ko
/var/crash/*
/boot/efi
/cm/images
/cm/shared/*
/proc/*
/sys/*
/tftpboot/*
/.autorelabel
/var/lib/logrotate.status
/var/lib/sss/*
/var/lib/systemd/random-seed
/var/spool/anacron/*
/cm/images/k8s-worker-dgx-h200-image-071525/boot/*
/cm/images/k8s-worker-dgx-h200-image-041725/boot/*
/cm/images/k8s-control-image-031325/*
/cm/node-installer/*
/cm/images/k8s-worker-dgx-h200-image-030725/boot/*
/cm/images/k8s-runai-control-image-001-031325/boot/*
/cm/images/k8s-worker-dgx-h200-image-kdump/*
/cm/images/test-image-13032025-3/boot/*
/cm/images/default-image/boot/*
/cm/images/k8s-runai-control-image-001/boot/*
/cm/images/k8s-worker-dgx-h200-1082-gold-image/boot/*
/cm/images/test-image-13032025-2/*
/cm/images/test-image-13032025-1/*
/cm/images/k8s-control-image-031325/boot/*
/cm/images/k8s-runai-control-image-051625/boot/*
/tftpboot/*
/cm/images/k8s-worker-dgx-h200-image-new-kernel-gpu-driver/*
/cm/images/k8s-worker-dgx-h200-image-050725/boot/*
/cm/images/k8s-worker-dgx-h200-image-031025-1/boot/*
/cm/images/k8s-worker-dgx-h200-image-001/*
/var/spool/cmd/monitoring/*
/cm/images/k8s-worker-dgx-h200-image-031025-1/*
/cm/images/login-image-1/boot/*
/cm/images/k8s-runai-control-image-051225/*
/cm/images/k8s-worker-dgx-h200-image-060525/boot/*
/cm/images/k8s-worker-dgx-h200-image-kdump/boot/*
/cm/images/login-image-1/*
/cm/images/k8s-worker-dgx-h200-1082-lustre-driver-added/*
/cm/images/k8s-worker-dgx-h200-image-041725/*
/cm/images/k8s-worker-dgx-h200-image-031025/boot/*
/cm/images/k8s-runai-control-image-051325/*
/cm/images/k8s-control-image-041725/*
/cm/images/test-image-13032025-clone/*
/cm/images/test-image-13032025-3/*
/cm/images/default-image/*
/cm/images/test-image-13032025-clone/boot/*
/cm/images/k8s-control-image-051225/*
/cm/conf/all/*
/cm/images/k8s-worker-dgx-h200-image-051625/*
/cm/images/k8s-worker-dgx-h200-image-070925/boot/*
/cm/images/k8s-runai-control-image-001/*
/cm/images/k8s-runai-control-image-051225/boot/*
/cm/images/test-image-13032025-1/boot/*
/cm/images/k8s-worker-dgx-h200-image-031025/*
/cm/images/k8s-worker-dgx-h200-image-001/boot/*
/cm/images/login-image-2/*
/cm/images/k8s-worker-dgx-h200-image-070925/*
/cm/images/k8s-control-image-041725/boot/*
/cm/images/login-image/boot/*
/cm/images/k8s-runai-control-image-041725/boot/*
/cm/images/k8s-worker-dgx-h200-image-051025/boot/*
/cm/images/k8s-runai-control-image-001-031325/*
/cm/images/k8s-runai-control-image-051625/*
/cm/images/k8s-control-image-001/*
/cm/images/k8s-worker-dgx-h200-image-071525/*
/cm/images/login-image/*
/cm/images/k8s-worker-dgx-h200-dgx-os7-image-062425/*
/cm/shared/*
/cm/images/k8s-control-image-051225/boot/*
/cm/images/k8s-worker-dgx-h200-image-060525/*
/cm/images/k8s-worker-dgx-h200-image-051025/*
/cm/images/k8s-worker-dgx-h200-1082-lustre-driver-added/boot/*
/cm/images/k8s-control-image-001/boot/*
/cm/images/dgx-os-6.3-h100-image/boot/*
/cm/images/k8s-worker-dgx-h200-dgx-os7-image-062425/boot/*
/cm/images/k8s-worker-dgx-h200-1082-gold-image/*
/cm/images/k8s-worker-dgx-h200-image-050725/*
/cm/images/k8s-control-image-052025/boot/*
/cm/images/login-image-2/boot/*
/cm/images/k8s-worker-dgx-h200-image-030725/*
/cm/images/k8s-worker-dgx-h200-image-new-kernel-gpu-driver/boot/*
/cm/images/dgx-os-6.3-h100-image/*
/cm/images/k8s-worker-dgx-h200-image-051625/boot/*
/cm/images/k8s-runai-control-image-051325/boot/*
/cm/images/test-image-13032025/*
/cm/images/k8s-runai-control-image-041725/*
/cm/images/test-image-13032025-2/boot/*
/cm/images/k8s-control-image-052025/*
/cm/images/test-image-13032025/boot/*
/var/*/rsyslog/imjournal.state
/etc/machine-id
/var/lib/logrotate/logrotate.status
/var/lib/systemd/timers/stamp-*
/var/lib/containerd
/cm/local/apps/containerd/var/etc/conf.d
/etc/systemd/system/multi-user.target.wants/containerd.service
/cm/local/apps/containerd/var/etc/conf.d/cri.toml
/cm/local/apps/containerd/var/etc/conf.d/nvidia-cri.toml
/cm/local/apps/containerd/var/etc/conf.d/custom-sandbox-image.toml
/etc/nginx/nginx.conf
/root/.kube
/root/.kube/*
/root/.byoh
/root/.byoh/*
/var/lib/kubelet/device-plugins
/var/lib/kubelet/device-plugins/*
/etc/kubernetes
/etc/kubernetes/*
/opt/containerd
/opt/containerd/*
/usr/libexec/kubernetes
/usr/libexec/kubernetes/*
/var/lib/kubelet
/var/lib/kubelet/*
/var/lib/cni
/var/lib/cni/*
/etc/cm-kube-version.json
/var/lib/etcd
/var/lib/etcd/*
/usr/lib/systemd/system/kubelet.service.d
/usr/lib/systemd/system/kubelet.service.d/*
/etc/cni
/opt/cni
/var/lib/calico
/run/calico
/etc/cni
/run/flannel/*
